# MongoDB/mongodump-backup.timer
[Unit]
Description=MongoDB Daily Backup

[Timer]
OnCalendar=daily
Persistent=true
Unit=mongodump-backup.service

[Install]
WantedBy=timers.target
